# crontab/crontab.lst
# Motiviz Cron

SHELL=/bin/bash

0 9 * * 1-5 source /home/dlink/.profile && export VCONF=/home/dlink/motiviz/conf/dev.yaml && /home/dlink/motiviz/lib/hip.py inspire 2>> /var/log/motiviz/inspire_error.log >> /var/log/motiviz/inspire.log
